# ipython/service/ipython/run
#!/bin/sh
exec 2>&1
cd /root/notebook
exec /usr/local/bin/ipython notebook --profile=nbserver --ip='*'

